# source/tb321fu-camera-rootfs-overlay/rootfs-overlay/etc/systemd/user/pipewire.service.d/60-y700-libcamera-paths.conf
[Service]
Environment=LD_LIBRARY_PATH=/opt/libcamera-y700/lib/aarch64-linux-gnu:/opt/libcamera-y700/lib
Environment=LIBCAMERA_IPA_PROXY_PATH=/opt/libcamera-y700/libexec/libcamera
Environment=LIBCAMERA_IPA_MODULE_PATH=/opt/libcamera-y700/lib/aarch64-linux-gnu/libcamera/ipa
Environment=LIBCAMERA_IPA_CONFIG_PATH=/opt/libcamera-y700/share/libcamera/ipa
